# scripts/smbits-autostart.service
[Unit]
Description=SMBITS ERP Portal (Flask App)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/erpnext-custom-web-app
ExecStart=%h/erpnext-custom-web-app/venv/bin/python %h/erpnext-custom-web-app/app.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
